# com.cf.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cf.bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
        <string>cf_bot.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/zonedout/discord/cf colour bot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/zonedout/discord/cf colour bot/bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zonedout/discord/cf colour bot/bot_error.log</string>
</dict>
</plist>
